# supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700
username=dummy
password=dummy

[supervisord]
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
username=dummy
password=dummy

[program:michalsoida]
command=/usr/local/bin/gunicorn -k gevent michalsoida:app -b 127.0.0.1:8000 -w 3 --log-file=-
directory=/app
user=root
autostart=true
autorestart=true
redirect_stderr=True
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
